# conf/supervisor/talebook.conf
[group:talebook]
programs=nginx,bootstrap,tornado

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=unexpected

; 启动自检：权限校验、nginx 配置检查、syncdb/migrate/update-config，逐项写入
; status.json；全部通过后才会 `supervisorctl start tornado`。失败时不重试
; （autorestart=false），nginx 与 /status.json、/status_page.html 仍然可访问。
[program:bootstrap]
command=python3 /var/www/talebook/webserver/self_check.py
directory=/var/www/talebook
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

; 不再随 group 自动启动，由 bootstrap 自检全部通过后拉起
[program:tornado]
command=gosu talebook:talebook python3 server.py --with-library=/data/books/library --port=8000 --host=127.0.0.1 --logging=debug --log-file-prefix=/data/log/talebook.log
directory=/var/www/talebook
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

